# examples/systemd/3gpp-mcp-update.service
[Unit]
Description=Update 3GPP specifications to latest versions
After=network.target

[Service]
Type=oneshot
User=3gpp-mcp
Group=3gpp-mcp
WorkingDirectory=/var/lib/3gpp-mcp
ExecStart=/usr/local/bin/3gpp-mcp update --db /var/lib/3gpp-mcp/3gpp.db --convert-doc
